# overlay/etc/systemd/system/set_time.service
[Unit]
Description=set time
After=multi-user.target getty.target

[Service]
Type=idle
ExecStart=/sbin/hwclock -s
[Install]
WantedBy=startup.target
