# overlay/etc/init.d/droid4-wlan
#!/sbin/openrc-run

# rc-update add droid4-wlan boot

mac=""

depend() {
	need dev
	need localmount
	after modules
	before wpa_supplicant
}

#
# Note that this does not consider the calibration data right now
#
read_mac() {
        tmp=$(busybox mktemp -d)
	if [ "${tmp}" = "" ]; then
		echo "Could not make a tmp directory"
		return
	fi

        if ! mount -o ro /dev/mmcblk1p7 /${tmp}; then
                echo "Could not mount pds"
                return
        fi

        if [ ! -f ${tmp}/wifi/nvs_map.bin ]; then
                echo "Could not find nvs_map.bin"
                umount ${tmp}
                return
        fi

        bytes=$(hexdump -C ${tmp}/wifi/nvs_map.bin | head -n1 | \
                        cut -d' ' -f6,7,8,9,14,15)

        for byte in ${bytes}; do
                if [ "${mac}" == "" ]; then
                        mac=${byte}
                else
                        mac=${byte}:${mac}
                fi
        done

        echo "Read WLAN MAC from nvs_map.bin: ${mac}"

        umount ${tmp}
}

start() {
	read_mac
	if [ "${mac}" != "" ]; then
		ifconfig wlan0 down
		ifconfig wlan0 hw ether ${mac}
	else
		echo "Could not read WLAN MAC"
		return 1
	fi
}
