# overlay/etc/init.d/droid4-buttons-handle
#!/sbin/openrc-run

depend() {
        after modules
}
export PYTHONPATH="/usr/local/share/droid4/python3_packages"
command="/usr/local/bin/buttons_handler"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
