# overlay/etc/init.d/droid4-hardware-demon
#!/sbin/openrc-run

depend() {
        after modules
	after droid4-startup
}
export PYTHONPATH="/usr/local/share/python"
command="/usr/local/bin/droid4-hardware-demon"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/root/.droid4/modem/logs/modem_demon.log"
