# overlay/etc/init.d/droid4-buttons-hendle
#!/sbin/openrc-run

depend() {
        after modules
}

command="/usr/local/bin/hardware/buttons_hendler.py"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
