# overlay/etc/init.d/droid4-startup
#!/sbin/openrc-run

depend() {
	need dev
}
command="/usr/local/bin/droid4-startup.sh"
command_args='start'
pidfile="/run/${RC_SVCNAME}.pid"
