# overlay/etc/profile.d/set_python_path.sh
export PYTHONPATH="/usr/local/share/python"
